# cron-s6/s6-overlay/s6-rc.d/cron-log/run
#!/bin/sh
exec echo dude
#exec logutil-service /var/log/cron
