# cron-s6/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv sh
env >> /etc/environment
#fdmove 2 1
sleep 1
exec 2>&1
cron
exec tail -f /var/log/syslogd/cron/current /var/log/cron.log
